# templates/rc_script.OpenBSD.j2
#!/bin/sh

# This file was generated by Ansible for fw3.home
# Do NOT modify this file by hand!

daemon="/usr/local/sbin/openvpn"
config_file="{{openvpn_system_etcdir }}/{{ __openvpn_service_name }}/{{ openvpn_default_instance_name }}.conf"
pid_file="/var/run/{{ __openvpn_service_name }}.pid"
dir="{{openvpn_system_etcdir }}/{{ __openvpn_service_name }}"
daemon_flags="--daemon {{ __openvpn_service_name }} --writepid $pid_file --cd $dir --config $config_file"

. /etc/rc.d/rc.subr

rc_reload (){
        if [ -f $pid_file ]; then
                kill -USR1 `cat $pid_file`
        else
                return 1
        fi
}


rc_cmd $1
